# supervisord/natmeter.ini
[program:natmeter]
priority=10
directory=/opt/django/natmeter/stun
command=/opt/django/natmeter/venv/bin/python manage.py crontab
user=www-data
autostart=true
autorestart=true
